# deploy/edge/nic-tune.service
[Unit]
Description=eSPX ingress NIC RX ring and IRQ tuning
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/espx/edge-nic-tune.env
ExecStart=/usr/local/bin/edge_nic_tune.sh apply
# RX ring and IRQ affinity are lost on driver reload/reboot; re-apply at boot.

[Install]
WantedBy=multi-user.target
