# deploy/cron/log-evacuate.cron
# log-evacuate.cron: long-lived S3 evacuator for rotated tracker segments.
# Prefer: docker compose --profile tools up -d log-evacuator
# Legacy one-shot restart every 5 minutes (not recommended for fsnotify watcher):

*/5 * * * * cd /opt/espx && docker compose --profile tools up -d log-evacuator >> /var/log/espx-evacuate.log 2>&1
