# deploy/cron/log-evacuate.cron
# log-evacuate.cron: rotate and ship tracker logs to the broker before disks fill.
# Runs every 5 minutes via docker compose tools profile.

*/5 * * * * cd /opt/espx && docker compose --profile tools run --rm log-evacuator >> /var/log/espx-evacuate.log 2>&1
